# upstart/ok-webfront.conf
# start ok-frontend node application - production
start on runlevel [2345]
stop on runlevel [016]