# upstart/ok-webfront-web-3000.conf
## Starts a node.js server using express.js for ok-webfront on port 3000
##
start on starting ok-webfront
stop on stopping ok-webfront

respawn

env NODE_ENV=production

pre-start script
  APP_DIR=/home/zohar/workspace/ok-webfront # change to install path of application
  RUN_AS=zohar                              # change to user the app should run under
  mkdir -p $APP_DIR/tmp/pids
  mkdir -p $APP_DIR/log
  chown -R $RUN_AS $APP_DIR/tmp/pids
  chown -R $RUN_AS $APP_DIR/log
end script

script
  APP_DIR=/home/zohar/workspace/ok-webfront # change to install path of application
  RUN_AS=zohar                              # change to user the app should run under
  PORT=3000                                 # change to port the express.js server should listen to
  PID_FILE="$APP_DIR/tmp/pids/express-$PORT.pid"

  ## Run via start-stop-daemon
  exec start-stop-daemon --start --chuid $RUN_AS \
                         --make-pidfile --pidfile $PID_FILE \
                         --chdir $APP_DIR \
                         --exec /usr/local/bin/node -- \
                         $APP_DIR/app.js --port=$PORT
end script
